# supervisord.conf
[program:diebotdie]
command=python run_diebotdie.py rules.json --port 6379 --host localhost --db 0
